# files/tmpfs/cache-ramdisk.mount
[Unit]
Before=local-fs.target

[Mount]
What=tmpfs
Where=<HOME>/.cache
Type=tmpfs
Options=noatime,nodiratime,nodev,nosuid,uid=1000,gid=1000,size=300M
