# files/chrome-ramdisk/chrome-ramdisk.mount
[Unit]
Before=local-fs.target

[Mount]
What=tmpfs
Where=<HOME>/.config/chromium
Type=tmpfs
Options=noatime,nodiratime,nodev,nosuid,uid=<UID>,gid=<GID>,size=300M
